# deploy/javasimple.service
[Unit]
Description=My Java Gradle App (Simple Service)
After=network.target

[Service]
Type=simple
User=appuser
Group=appuser
WorkingDirectory=/opt/myapp

# Use full path — no ${} in systemd allowed!
Environment=JAVA_HOME=/usr/lib/jvm/java-21-openjdk-amd64
ExecStart=/usr/lib/jvm/java-21-openjdk-amd64/bin/java -jar /opt/myapp/app.jar

SuccessExitStatus=143
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target


