# deploy/javaforking.service

[Unit]
Description=My Java Gradle App (Forking Service)
After=network.target

[Service]
Type=forking
User=appuser
Group=appuser
WorkingDirectory=/opt/myapp
ExecStart=/opt/myapp/start-wrapper.sh
ExecStop=/bin/kill -15 $MAINPID
SuccessExitStatus=143
Restart=on-failure

[Install]
WantedBy=multi-user.target

